# packages/server/scripts/launchd/com.zhixuan92.mmagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zhixuan92.mmagent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/mmagent</string>
    <string>serve</string>
    <string>--http</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>PLACEHOLDER_HOME/.multi-model/logs/daemon.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>PLACEHOLDER_HOME/.multi-model/logs/daemon.stderr.log</string>
</dict>
</plist>
