# packages/server/scripts/systemd/mma.service
[Unit]
Description=multi-model-agent HTTP daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mma serve --http
Restart=always
RestartSec=5
StandardOutput=append:%h/.mma/logs/daemon.stdout.log
StandardError=append:%h/.mma/logs/daemon.stderr.log

[Install]
WantedBy=default.target
